# bin/tunel_server.service
[Unit]
Description=TunelMxServer

[Service]
ExecStart=$HOME/.tunelmx/tunel_server

[Install]
WantedBy=multi-user.target
