# bin/tunel_client.service
[Unit]
Description=TunelMxClient

[Service]
ExecStart=$HOME/.tunelmx/tunel_client

[Install]
WantedBy=multi-user.target
